# etc/supervisord.d/c0.ini
[program:carbon-cache-c0]
command = /usr/bin/python /bin/carbon-cache --instance=c0 --conf=/etc/carbon/c0.conf --pidfile=/var/run/carbon-cache_c0.pid --debug start
redirect_stderr=true
stdout_syslog = true
autorestart = true
