# etc/supervisord.d/c1.ini
[program:carbon-cache-c1]
command = /usr/bin/python /bin/carbon-cache --instance=c1 --conf=/etc/carbon/c1.conf --pidfile=/var/run/carbon-cache_c1.pid --debug start
redirect_stderr=true
stdout_syslog = true
autorestart = true
