# etc/supervisord.d/etcd.ini
[program:etcd]
directory=/var/lib/etcd
command=/root/bin/start_etcd.sh
autostart=false
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log

