# etc/supervisord.d/graphite-api.ini
[program:graphite-api]
command = /usr/bin/gunicorn -w2 graphite_api.app:app -b 127.0.0.1:8888
redirect_stderr=true
stdout_syslog = true
autorestart = true
