# etc/supervisord.d/kibana.ini
[program:kibana]
command = /opt/kibana/bin/kibana
redirect_stderr=true
stdout_logfile = syslog
