# etc/supervisord.d/nginx_gapi.ini
[program:nginx_api]
command = /usr/sbin/nginx -c /etc/nginx_gapi/nginx.conf
redirect_stderr=true
stdout_syslog=true
