# etc/supervisord.d/r0.ini
[program:carbon-relay-r0]
command = /usr/bin/python /bin/carbon-relay --instance=r0 --conf=/etc/carbon/r0.conf --pidfile=/var/run/carbon-relay_r0.pid --debug start
redirect_stderr=true
stdout_syslog = true
autorestart = true
